# docker/rootfs/etc/services.d/nginx/run
#!/usr/bin/execlineb -P
# Learn more detail about `execlineb` at https://danyspin97.org/blog/getting-started-with-execline-scripting/

/usr/sbin/nginx