# docker/rootfs/etc/services.d/crond/run
#!/usr/bin/with-contenv sh
#
# Run busybox to manage crontab for non-root users
#

/usr/sbin/crond -f -d 7
