# docker/rootfs/etc/services.d/php5-fpm/run
#!/usr/bin/with-contenv sh

/usr/bin/php-fpm5